# supervisor/conf.d/dump1090.conf
[program:dump1090]
priority=20
user=docker
group=nogroup
command=/usr/local/bin/dump1090 --net --write-json /var/lib/dump1090/public_html/data --quiet --gain 35
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
