# supervisor/conf.d/python_http.conf
[program:python_http]
priority=10
command=/usr/bin/python -m SimpleHTTPServer 8080 
user=docker
group=nogroup
directory=/var/lib/dump1090/public_html
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
